# deploy/supervisord.conf
# Runs the two agents + nginx in one container. Children inherit the container
# environment (ANTHROPIC_API_KEY, GMI_MAAS_*, etc.) and add per-process vars.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# Resume agent (Python). Runs as non-root: the claude CLI refuses to run as root.
[program:resume]
command=python agent1.py
directory=/app
user=app
environment=AGENT_PORT="8001",AGENT_DB="/data/agent1.db",HOME="/home/app"
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Interview agent (Node). Reads the resume agent's agent1.db read-only via AGENT_DB.
[program:interview]
command=node --enable-source-maps dist/index.js
directory=/app/interview
user=app
environment=PORT="8002",INTERVIEW_DB="/data/interview.db",AGENT_DB="/data/agent1.db",INTERVIEW_ARTIFACT_ROOT="/data/interview-artifacts",HOME="/home/app"
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Reverse proxy on 8080 (the GMI web port). Runs as root.
[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
